# .config/systemd/user/redshift.service
[Unit]
Description=Redshift

[Service]
ExecStart=/usr/bin/redshift-gtk

[Install]
WantedBy=user-applications.target

