# .config/systemd/user/slack.service
[Unit]
Description=Slack Desktop

[Service]
ExecStart=/usr/bin/slack

[Install]
WantedBy=user-applications.target

